# scripts/com.cryptobot.news-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Standalone news scraper — runs every 15 minutes via launchd -->
    <!-- Install: cp scripts/com.cryptobot.news-scraper.plist ~/Library/LaunchAgents/ -->
    <!-- Load:    launchctl load ~/Library/LaunchAgents/com.cryptobot.news-scraper.plist -->
    <!-- Unload:  launchctl unload ~/Library/LaunchAgents/com.cryptobot.news-scraper.plist -->

    <key>Label</key>
    <string>com.cryptobot.news-scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/daniil/Projects/crypto-investment-bot/scripts/run_chrome_mcp_scraper.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/daniil/Projects/crypto-investment-bot</string>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/daniil/Projects/crypto-investment-bot/data/news-scraper.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/daniil/Projects/crypto-investment-bot/data/news-scraper.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
